# alc298-hda-init.service
[Unit]
Description=Initialize ALC298 codec
After=sound.target
Wants=sound.target

[Service]
Type=oneshot
ExecStart=/usr/share/alc298-hda-init/alc298-hda-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
